# Algoritmos/SERVICIOS-AURALIS/auralis-subscriber final/auralis-subscriber.service
[Unit]
Description=Auralis MQTT Subscriber
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/default/auralis-subscriber
User=auralis
Group=auralis
WorkingDirectory=/opt/auralis-subscriber
ExecStart=/opt/auralis-subscriber/.venv/bin/python /opt/auralis-subscriber/run.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
